# auto_upgrade.timer
[Unit]
Description=Run auto_upgrade.py daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
